# notifications/docker/build/supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/supervisor/supervisord.log

[program:async_user_created]
command=php bin/console messenger:consume -vv async_user_created --memory-limit=128M --time-limit=3600
autostart=true
autorestart=true
stderr_logfile=/var/log/queue-err.log
stdout_logfile=var/log/queue.log
process_name=%(program_name)s_%(process_num)02d
